# roles/worker/files/poll.service
[Unit]
Description=worker

[Service]
ExecStart=/usr/bin/java /worker/src/main/java/Worker.java

[Install]
WantedBy=multi-user.target