# roles/result/files/result.service
[Unit]
Description=result

[Service]
ExecStart=/usr/bin/node /home/result/server.js

[Install]
WantedBy=multi-user.target