# roles/poll/files/poll.service
[Unit]
Description=poll

[Service]
ExecStart=/usr/bin/python /home/poll/app.py

[Install]
WantedBy=multi-user.target